# snapcast-server/rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting dbus-daemon..."

# https://gitlab.freedesktop.org/dbus/dbus/-/issues/441#note_1721753
ulimit -n 1024

exec dbus-daemon --system --nofork
