# snapcast-server/rootfs/etc/services.d/gmediarender/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Starting DLNA Renderer' 
device_name=$(bashio::config 'renderer_name')

exec /usr/local/bin/gmediarender -f ${renderer_name} -p 50000 --gstout-audiopipe 'audioresample ! audio/x-raw, rate=48000, format=S16LE ! filesink location=/share/snapfifo'